# other/terminal-sharing-start.timer
[Unit]
Description=terminal sharing start every 15 seconds

[Timer]
OnBootSec=15s
OnUnitInactiveSec=15s
AccuracySec=1s

[Install]
WantedBy=timers.target
